# roles/app/files/compress-stepup-authnlogs.cron
#!/bin/bash

# Compress stepup authentication logs older than 1 day
find /var/log/ -maxdepth 1 -name "stepup-authentication-*.log" -mtime +1 -exec bzip2 -9 {} \;