# roles/app/files/delete-old-stepup-authnlogs.cron
#!/bin/bash

# delete authentication logs older than 31 days
find /var/log/ -maxdepth 1 -name "stepup-authentication-*.log.bz2" -mtime +31 -exec rm {} \;